# scripts/crontab
0   3   *   *   0   %%USER_NAME%%    cd %%BIN_DIR%%; ./ingest.sh
