# scripts/crontab
0   1   *   *   6   %%USER_NAME%%    cd %%BIN_DIR%%; ./ingest.sh
