# scripts/crontab
0   1   *   *   6   root    cd ingest-dir; ./ingest.sh
